# media-automount@.service
[Service]
Type=forking
GuessMainPID=no
ExecStart=/bin/media-automount %I
